# nfpm/systemd/prometheus-nagios-exporter.service
[Unit]
Description=Nagios Exporter
Documentation=https://github.com/wbollock/nagios_exporter
Wants=network-online.target nagios.service
After=network-online.target nagios.service

[Service]
EnvironmentFile=-/etc/default/prometheus-nagios-exporter
User=root
Group=root
Type=simple
ExecStart=/usr/local/bin/prometheus-nagios-exporter $ARGS

[Install]
WantedBy=multi-user.target
